# meshtastic-bbs.service
[Unit]
Description=Meshtastic PTT-BBS 伺服器
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/meshtastic-bbs
ExecStart=/usr/bin/python3 /home/pi/meshtastic-bbs/meshtastic_bbs_server.py --port /dev/ttyUSB0
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
